# .crontab
45 23 * * * bash -c "CLUSTER=\$(echo \$(hostname) | cut -d '.' -f 2); if [[ \"gautschi\" == *\"\$CLUSTER\"* ]]; then SLURM_CONFIG_PATH=/home/\$(whoami)/rcac-utils; else CLUSTER=\$(echo \$(hostname) | cut -d'.' -f 1); SLURM_CONFIG_PATH=/home/\$CLUSTER/a/\$(whoami)/rcac-utils; fi; bash \$SLURM_CONFIG_PATH/auto_env_export.bash"
